# openqa/recover-nfs.timer
[Unit]
Description=Automatically recover stall nfs shares.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
