# openqa/salt-minion-restarter.service

[Unit]
Description=Restart salt minion service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart salt-minion
